# files/ckan.cron
@reboot root /usr/local/bin/reset_ckan_pass.sh